# rust-watchdog.service
[Unit]
Description=Rust LinuxGSM watchdog (update + oxide + restart)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=rustserver
Group=rustserver
WorkingDirectory=/home/rustserver
ExecStart=/usr/bin/python3 /home/rustserver/rust_watchdog.py --config /home/rustserver/rust_watchdog.json
Restart=always
RestartSec=5
KillMode=process

# Hardening (optional)
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
