# rust-watchdog.service
[Unit]
Description=Rust LinuxGSM watchdog (update + mods-update + restart)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=rustserver
Group=rustserver
WorkingDirectory=/home/rustserver/rust-linuxgsm-watchdog

# Create a private runtime dir under /run for a lockfile
RuntimeDirectory=rust-watchdog
RuntimeDirectoryMode=0755

# HARD LOCK: prevent multiple watchdog instances (and duplicate RustDedicated)
# venv method:
ExecStart=/usr/bin/flock -n /run/rust-watchdog/lock \
  /home/rustserver/rust-linuxgsm-watchdog/.venv/bin/python \
  /home/rustserver/rust-linuxgsm-watchdog/rust_watchdog.py \
  --config /home/rustserver/rust-linuxgsm-watchdog/rust_watchdog.json

# alt method (no venv):
# ExecStart=/usr/bin/flock -n /run/rust-watchdog/lock \
#   /usr/bin/python3 \
#   /home/rustserver/rust-linuxgsm-watchdog/rust_watchdog.py \
#   --config /home/rustserver/rust-linuxgsm-watchdog/rust_watchdog.json

# Restart only when the watchdog crashes / exits non-zero
Restart=on-failure
RestartSec=5

# On stop/restart, kill the whole cgroup (watchdog + tmux/RustDedicated it spawned).
# (KillMode defaults to control-group; leaving it unset is fine.)
KillSignal=SIGINT
TimeoutStopSec=30

NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
